# systemd_units/tuner.service
[Unit]
Description=Tuner
Documentation=https://github.com/adnanh/webhook

[Service]
ExecStart=/usr/bin/ruby /opt/tuner.rb
