# systemd_units/tuner.service
[Unit]
Description=Tuner

[Service]
ExecStart=/usr/bin/ruby /opt/tuner.rb
